# autostart/launch.sh.desktop
[Desktop Entry]
Exec=/home/regal/.config/polybar/colorblocks/launch.sh
Icon=
Name=launch.sh
Path=
Terminal=False
Type=Application
